# templates/temp.service.j2
[Unit]
Description=Our example of go application
After=network.target

[Service]
Type=simple
ExecStart=/home/vagrant/exam
PIDFile=/var/run/simple

[Install]
WantedBy=default.target
